# manager_supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=root                 ; (default is current user, required if root)
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


[program:webui]
command=/root/anaconda3/bin/python run_manager_web.py
directory=/root/360-manager
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/root/360-manager/logs/webui.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info


[program:celery_task]
command=/root/anaconda3/bin/celery -A backend.celery_task worker --loglevel=INFO  --concurrency 4
directory=/root/360-manager
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/root/360-manager/logs/celery_task.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info

[program:celery_flower]
command=/root/anaconda3/bin/celery flower -A backend.celery_task --address=0.0.0.0 --port=5658
directory=/root/360-manager
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/root/360-manager/logs/celery_task.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info